# templates/firedancer-docker.service
[Unit]
Description=Firedancer in docker
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory={{ firedancer_workdir }}
ExecStartPre=/usr/bin/docker compose down --remove-orphans
ExecStart=/usr/bin/docker compose -f docker-compose.yml up
ExecStop=/usr/bin/docker compose down --remove-orphans

# Restart policy
Restart=always
TimeoutSec=infinity
RestartSec=5s

[Install]
WantedBy=multi-user.target